# tauri/dev.getminutes.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Minutes Voice Memo Watcher
         Watches ~/.minutes/inbox/ for new audio files and auto-processes them.
         Install: cp dev.getminutes.watcher.plist ~/Library/LaunchAgents/
         Load:    launchctl load ~/Library/LaunchAgents/dev.getminutes.watcher.plist
         Unload:  launchctl unload ~/Library/LaunchAgents/dev.getminutes.watcher.plist
         Logs:    tail -f ~/.minutes/logs/watcher-stdout.log
    -->
    <key>Label</key>
    <string>dev.getminutes.watcher</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update this path to your minutes binary -->
        <string>/Users/silverbook/.cargo/bin/minutes</string>
        <string>watch</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/silverbook/.minutes/logs/watcher-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/silverbook/.minutes/logs/watcher-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/silverbook</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
